# extra/lib/systemd/system/keymapperd.service
[Unit]
Description=Keymapper Daemon

[Service]
ExecStart=keymapperd

[Install]
WantedBy=multi-user.target
