# extra/lib/systemd/system/keymapperd.service
[Unit]
Description=Keymapper Daemon

[Service]
ExecStart=/bin/sh -c "keymapperd"

[Install]
WantedBy=multi-user.target
